# rpi-nitro/etc/systemd/system/led_mqtt.service
[Service]
ExecStart=/usr/bin/python /srv/mqtt/led_mqtt.py

[Install]
WantedBy=multi-user.target
