# rpi-nitro/etc/systemd/system/motion_ctrl_mqtt.service
[Service]
ExecStart=/usr/bin/python /srv/mqtt/motion_ctrl_mqtt.py
Requires=motion.service
After=motion.service

[Install]
WantedBy=multi-user.target
