# rpi-nitro/etc/systemd/system/tsl2591_mqtt.service
[Service]
ExecStart=/usr/bin/python /srv/mqtt/tsl2591_mqtt.py

[Install]
WantedBy=multi-user.target
