# rpi-piva/etc/systemd/system/htu21df_mqtt.service
[Service]
ExecStart=/usr/bin/python /srv/mqtt/htu21df_mqtt.py

[Install]
WantedBy=multi-user.target
