# rpi-piva/etc/systemd/system/pir_mqtt.service
[Service]
ExecStart=/usr/bin/python /srv/mqtt/pir_mqtt.py

[Install]
WantedBy=multi-user.target
