# rpi-ney/etc/systemd/system/ctrl_mqtt.service
[Service]
ExecStart=/usr/bin/python /srv/mqtt/ctrl_mqtt.py
Requires=mosquitto.service
After=mosquitto.service
[Install]
WantedBy=multi-user.target
