# rpi-ney/etc/systemd/system/node.service
[Service]
ExecStart=/usr/bin/node /srv/node/src/server.js
Requires=mysqld.service
After=mysqld.service
[Install]
WantedBy=multi-user.target